# .config/systemd/user/wallpapers.timer
[Unit]
Description=Change wallpaper every 10 minutes

[Timer]
OnCalendar=*:0/10
Persistent=true

[Install]
WantedBy=multi-user.target
